# contrib/patchup-agent-refresh.timer
[Unit]
Description=patchup-agent: refresh system state
PartOf=patchup-agent.service

[Timer]
OnActiveSec=15sec
OnUnitActiveSec=15min
