# contrib/systemd/patchup-agent.socket
[Unit]
PartOf=patchup-agent.service

[Socket]
ListenStream=/run/patchup.sock
SocketMode=0660
Accept=no
