# contrib/systemd/patchup-agent.socket
[Unit]
PartOf=patchup-agent.service

[Socket]
ListenStream=/run/patchup/socket
SocketMode=0660
DirectoryMode=0750
Accept=no
