# contrib/openrc/patchup-agent.cron
#!/bin/sh
# Check for pending system updates

/usr/bin/patchup agent --refresh /run/patchup/socket
exit 0
